# materials/class05/src/former.ini
[program:former]
command=uwsgi
    --plugins=python
    --http-socket=0.0.0.0:80
    --wsgi-file /opt/webcode/former/process/webrunner.py
    --static-map /form=/opt/webcode/former/form/index.html
    --processes=5
    --master
    --pidfile=/tmp/former.pid
    --vacuum
    --max-requests=5000
stopsignal=QUIT
autostart=true
startretries=10
startsecs=0
stopwaitsecs=10
stopasgroup=true

stdout_logfile=/var/log/webapps/former_stdout.log
stdout_logfile_maxbytes=60MB
stdout_logfile_backups=4
stdout_capture_maxbytes=4MB

stderr_logfile=/var/log/webapps/former_stderr.log
stderr_logfile_maxbytes=60MB
stderr_logfile_backups=4
stderr_capture_maxbytes=4MB
